# Library/LaunchAgents/com.llama-swap.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.llama-swap</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/llama-swap</string>
        <string>--config</string>
        <string>${HOME}/.llama/config.yaml</string>
        <string>--listen</string>
        <string>localhost:11435</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>LLAMA_SERVER</key>
        <string>/opt/homebrew/bin/llama-server</string>
        <key>LLAMA_MODELS_DIR</key>
        <string>${HOME}/.llama/models</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>${HOME}/.llama/llama-swap.log</string>
    <key>StandardErrorPath</key>
    <string>${HOME}/.llama/llama-swap.log</string>
</dict>
</plist>